# email-stack.service
[Unit]
Description=email-stack
Requires=docker.service
After=docker.service

[Service]
Restart=always

EnvironmentFile=/opt/email-stack/.env
ExecStartPre=-/usr/bin/docker-compose -f /opt/email-stack/docker-compose.yml down
ExecStartPre=-/opt/email-stack/email-stack-ctl.sh createvol
ExecStart=/usr/bin/docker-compose -f /opt/email-stack/docker-compose.yml up
ExecStop=/usr/bin/docker-compose -f /opt/email-stack/docker-compose.yml down

[Install]
WantedBy=multi-user.target
